# src/main/assembly/dist/install/dd-dataverse-authenticator.service
[Unit]
Description=Dd Dataverse Authenticator Service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-dataverse-authenticator/bin/dd-dataverse-authenticator server /etc/opt/dans.knaw.nl/dd-dataverse-authenticator/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143
Environment="JAVA_OPTS=-Xms64m"
Environment="LC_ALL=en_US.UTF-8"

User=dd-dataverse-authenticator
Group=dd-dataverse-authenticator

[Install]
WantedBy=multi-user.target
